# Systemd/podman_port443.socket
[Unit]
Description=Forward port 443 traffic to Unix socket

[Socket]
ListenStream=443

[Install]
WantedBy=sockets.target
